# supervisord-litecoin.conf
[program:litecoind]
command=/start-litecoind.sh
numprocs=1
autostart=true
autorestart=false
