# supervisord-stratum.conf
[program:twistd]
command=/start-stratum.sh
numprocs=1
autostart=true
autorestart=false
